# signalhive/deploy/supervisord.conf
[program:signalhive-web]
command=/opt/signalhive/venv/bin/gunicorn -c gunicorn.conf.py wsgi:app
directory=/opt/signalhive
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/signalhive/web.log
stderr_logfile=/var/log/signalhive/web_error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:signalhive-engine]
command=/opt/signalhive/venv/bin/python -m app.engine.runner
directory=/opt/signalhive
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/signalhive/engine.log
stderr_logfile=/var/log/signalhive/engine_error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:signalhive-collector-tg]
command=/opt/signalhive/venv/bin/python -m app.collectors.telegram_runner
directory=/opt/signalhive
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/signalhive/collector_tg.log
stderr_logfile=/var/log/signalhive/collector_tg_error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
